# debian/libpigpiod/postinst
#!/bin/sh
# postinst script for indi-asi-power

sudo systemctl daemon-reload
sudo systemctl enable pigpiod.service
sudo systemctl start pigpiod.service
