# debian/indi-rpi-gpio/postinst
#!/bin/sh
# postinst script for indi-rpi-gpio
sudo systemctl daemon-reload
sudo systemctl enable pigpiod.service
sudo systemctl start pigpiod.service
